# shared/stow/systemd/.config/systemd/user/polkit-kde-agent.service
[Unit]
Description=KDE PolicyKit authentication agent
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=/usr/lib/polkit-kde-authentication-agent-1
Restart=on-failure

[Install]
WantedBy=graphical-session.target
